# bin/supervisor.hippo.conf
[program:hippo]
command=python3 /opt/hippo/main.py
environment=DISPLAY=:1
priority=20
stdout_logfile=/var/log/hippo.log
redirect_stderr=true
startsecs=10
startretries=2